# go/BuildingRESTfulWebServiceswithGo/s1/myserver.ini
[program:myserver]
command=/home/ali/dev/code/graveyard-of-codes/go/BuildingRESTfulWebServiceswithGo/s1/bin/romainserver
autostart=true
autorestart=true
redirect_stderr=true
